# docker/celery_app/supervisor.conf
[unix_http_server]
file=/var/run/supervisor.sock                       ; path to your socket file

[supervisord]
nodaemon=true
childlogdir=/var/log/supervisord/
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock         ; use a unix:// URL  for a unix socket


[program:worker-app]
command=celery -A worker_app worker --loglevel=info
stdout_logfile=/var/log/worker-app.log
autostart=true
autorestart=true
user=www-data
stopsignal=KILL
numprocs=1



[program:worker-beats]
command=celery -A worker_app beat --loglevel=info
stdout_logfile=/var/log/worker-beats.log
autostart=true
autorestart=true
user=www-data
stopsignal=KILL
numprocs=1

